# tools/ops/systemd/vicuna-policy-rollout.timer
[Unit]
Description=Vicuña Policy Live Rollout Controller Timer

[Timer]
OnCalendar=hourly
Persistent=true
Unit=vicuna-policy-rollout.service

[Install]
WantedBy=timers.target
