# tools/ops/systemd/vicuna-runpod-capture-sync.timer
[Unit]
Description=Vicuña RunPod Experimental Capture Sync Timer

[Timer]
OnBootSec=30s
OnUnitActiveSec=15s
Unit=vicuna-runpod-capture-sync.service

[Install]
WantedBy=timers.target
