# tools/ops/systemd/vicuna-runpod-mistral-relay.system.service
[Unit]
Description=Vicuña RunPod Mistral Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@VICUNA_SERVICE_USER@
Group=@VICUNA_SERVICE_GROUP@
WorkingDirectory=@VICUNA_REPO_ROOT@
EnvironmentFile=-@VICUNA_ENV_FILE@
ExecStart=@VICUNA_REPO_ROOT@/tools/ops/run-with-service-logs.sh runpod-mistral-relay @VICUNA_REPO_ROOT@/tools/ops/run-runpod-mistral-relay.sh
Restart=always
RestartSec=2
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
