# tools/ops/systemd/vicuna-runpod-capture-sync.system.service
[Unit]
Description=Vicuña RunPod Experimental Capture Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=@VICUNA_REPO_ROOT@
EnvironmentFile=-@VICUNA_ENV_FILE@
ExecStart=@VICUNA_REPO_ROOT@/tools/ops/run-with-service-logs.sh runpod-capture-sync @VICUNA_REPO_ROOT@/tools/ops/runpod-sync-experimental-capture.sh
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
